# SystemdFiles/autoupdate.service
[Unit]
Description = Arch Auto Update Service

[Service]
Type=oneshot
ExecStart=/home/htpc/autoupdate/autoupdate.sh
